# idmtools_platform_local/docker_scripts/services/nginx/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid idmtools nginx -c /app/nginx.conf