# idmtools_platform_local/docker_scripts/services/nginx/run
#!/usr/bin/with-contenv sh

echo Starting UI
nginx -g 'daemon off;' -c /app/nginx.conf