# idmtools_platform_local/docker_scripts/services/uwsgi/run
#!/usr/bin/with-contenv sh
cd /app
exec s6-setuidgid idmtools uwsgi --ini /app/uwsgi.ini